# services/mdns-publisher/mdns-publisher.service
[Unit]
Description=Avahi CNAME/mDNS publisher

[Service]
ExecStart=mdns-publisher.sh

[Install]
WantedBy=multi-user.target
